# rootfs/etc/inittab
# Системная инициализация
::sysinit:/etc/init.d/rcS
ttyS0::askfirst:/bin/sh
# Консольные терминалы
tty1::respawn:/sbin/getty -L tty1 0 vt100
tty2::respawn:/sbin/getty -L tty2 0 vt100
tty3::respawn:/sbin/getty -L tty3 0 vt100
tty4::respawn:/sbin/getty -L tty4 0 vt100

# Перезагрузка при Ctrl+Alt+Del
::ctrlaltdel:/sbin/reboot

# Выключение системы
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
